# infra/com.career-ops.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.career-ops.healthcheck</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>cd $HOME/code/career-ops &amp;&amp; uv run python -m agents.tools.session_health</string>
    </array>
    <key>StartInterval</key>
    <integer>21600</integer>
    <key>StandardOutPath</key>
    <string>/Users/Shared/Library/Logs/career-ops/healthcheck.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/Library/Logs/career-ops/healthcheck-err.log</string>
</dict>
</plist>
